# scripts/uniacc-backend.service
[Unit]
Description=UNIACC Dashboard Backend API
After=network.target

[Service]
Type=simple
User=cl159906175
WorkingDirectory=/opt/Uniacc-dashboard-admision/backend
Environment="NODE_ENV=production"
EnvironmentFile=/opt/Uniacc-dashboard-admision/backend/.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StandardOutput=append:/opt/Uniacc-dashboard-admision/logs/backend.log
StandardError=append:/opt/Uniacc-dashboard-admision/logs/backend-error.log

[Install]
WantedBy=multi-user.target
